# public/infrastructure/public.service
[Unit]
Description=Nous Ergon Public Site (Streamlit)
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/alpha-engine-dashboard/public
ExecStart=/home/ec2-user/alpha-engine-dashboard/.venv/bin/streamlit run app.py \
          --server.port=8502 --server.address=127.0.0.1 \
          --server.headless=true --browser.gatherUsageStats=false
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
