# infrastructure/systemd/boot-pull.timer
[Unit]
Description=Daily code-pull at 12:00 UTC (5am PDT / 4am PST)
Requires=boot-pull.service

[Timer]
# 12:00 UTC every day. This is 5am PT during PDT (Mar-Nov) and 4am PT
# during PST (Nov-Mar). Chosen to run before Brian wakes up so any pull
# failure is visible in the morning and can be addressed before the
# Saturday pipeline fires at 17:00 PT (00:00 UTC Sat).
OnCalendar=*-*-* 12:00:00
# If the instance was stopped at 12:00 UTC (rare for the always-on
# micro, but possible during maintenance), run the next time systemd
# comes up instead of skipping the slot entirely.
Persistent=true
# Random 5-minute jitter so multiple instances in the same account
# don't hammer GitHub at exactly the same second (only relevant if
# more instances are added later — harmless for now).
RandomizedDelaySec=300
Unit=boot-pull.service

[Install]
WantedBy=timers.target
